# mctp_emulator/service_files/xyz.openbmc_project.mctp-emulator.service
[Unit]
Description=MCTP Emulator

[Service]
ExecStart=/usr/bin/mctp-emulator
SyslogIdentifier=mctp-emulator

[Install]
WantedBy=multi-user.target
